# cron/crontab
SHELL=/bin/zsh
*/10 * * * * find ~ -type f -name .DS_Store 2>/dev/null -delete
