# cron/crontab
*/10 * * * * root find ~ -type f -name .DS_Store 2>/dev/null -delete
0 5 * * * root source ~/.zshrc && refresh
